# docker/supervisord.conf
[supervisord]
nodaemon = true

[program:opensshd]
command = /usr/sbin/sshd -D
